# .devbox/ansible/roles/symfony_dev/files/devbox-watch.conf

[program:devbox-watch]

command=/home/vagrant/project/node_modules/.bin/gulp watch
directory=/home/vagrant/project
user=vagrant
autostart=false
environment=HOME="/home/vagrant"
